# src/systemd/apw.timer
[Unit]
Description=Arch Package Watcher (Hourly)
PartOf=apw.service

[Timer]
OnCalendar=hourly
Persistent=true
