# src/systemd/apw.service
[Unit]
Description=Arch Package Watcher Update

[Service]
Type=oneshot
ExecStart=/usr/bin/apw update
