# src/systemd/apw.timer
[Unit]
Description=Arch Package Watcher (Hourly)

[Timer]
OnCalendar=hourly
Persistent=true
